# systemd/inbox-watcher.service
[Unit]
Description=Yggdrasil-DEW Ecosistema - Inbox Watcher Daemon
Documentation=https://github.com/alvarofernandezmota-tech/yggdrasil-dew
After=network-online.target tailscaled.service
Wants=network-online.target

[Service]
Type=simple
User=varopc
Group=varopc
WorkingDirectory=/srv/yggdrasil-dew
ExecStart=/usr/local/bin/inbox-watcher.sh
Restart=on-failure
RestartSec=10s
StandardOutput=append:/var/log/inbox-watcher.log
StandardError=append:/var/log/inbox-watcher.log
Environment=HOME=/home/varopc
Environment=YGG_LOG_FILE=/var/log/yggdrasil-dew-maintenance.log
Environment=YGG_REPO_DIR=/srv/yggdrasil-dew

# Protección del sistema
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/srv/yggdrasil-dew /var/log
PrivateTmp=true

[Install]
WantedBy=multi-user.target
